# supervisord-activator.conf
[program:activator]
command=/start-activator.sh
numprocs=1
autostart=true
autorestart=true
